# backend/juci-core/etc/init.d/juci
#!/bin/sh /etc/rc.common

START=12

USE_PROCD=1
NAME=juci

start_service() {
	# todo: maybe add juci service here and start uhttpd etc. 
	# this will update index.html each time the box boots. 
	juci-update
}

stop() {
	return 0
}

reload() {
	juci-update
}

